# script/nezha-agent
#!/sbin/openrc-run
supervisor=supervise-daemon
name=nezha-agent
output_log=/var/log/${name}.log
error_log=/var/log/${name}.err
SERVER="nz_grpc_host:nz_grpc_port" #Dashboard 地址 ip:port
SECRET="nz_client_secret" #SECRET
pidfile="/run/${RC_SVCNAME}.pid"
command="/opt/nezha/agent/nezha-agent"
command_args="-s ${SERVER} -p ${SECRET}"
command_background=true

depend() {
	need net
}

start_pre() {
	checkpath -f -m 0644 -o root:root "/var/log/${name}.log"
}
