# script/nezha-dashboard
#!/sbin/openrc-run
supervisor=supervise-daemon
name=nezha-dashboard
output_log=/var/log/${name}.log
error_log=/var/log/${name}.err
pidfile="/run/${RC_SVCNAME}.pid"
command="/opt/nezha/dashboard/app"
command_args=""
command_background=true
directory="/opt/nezha/dashboard"

depend() {
	  need net
}

start_pre() {
    checkpath -f -m 0644 -o root:root "/var/log/${name}.log"
}
